# modules/syslinux/files/setboot-reg.service
[Unit]
Description=Set default boot to registration

[Service]
Type=oneshot
ExecStart=/usr/bin/cp /boot/syslinux/registration.cfg /boot/syslinux/syslinux.cfg

[Install]
WantedBy=multi-user.target
